# rootfs/etc/services.d/php-fpm7.0/run
#!/usr/bin/execlineb -P
cd /usr/sbin
php-fpm7.2 -F

