# rootfs/etc/services.d/refresh/run
#!/usr/bin/execlineb -P
cd /opt/scripts/
./refresh.sh

