# rootfs/etc/services.d/redis/run
#!/usr/bin/execlineb -P
cd /usr/bin
redis-server /opt/config/nzedb/redis.conf

