# rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bash
echo "Starting redis";
exec /usr/bin/redis-server /opt/config/nzedb/redis.conf

